# 11-midihub.rules
ACTION=="add|remove", SUBSYSTEM=="sound", ENV{ID_USB_INTERFACES}=="*010300*", RUN+="/usr/bin/python3 /home/oermens/midihub/midihub.py"
