# 11-midihub.rules
ACTION=="add|remove", SUBSYSTEM=="usb", ENV{ID_USB_INTERFACES}=="*010300*", RUN+="/bin/systemctl start midihub.service"
