# 33-midiusb.rules
# Trigger when a USB MIDI device is added or removed
ACTION=="add|remove", SUBSYSTEM=="sound", KERNEL=="card*", \
    RUN+="/usr/bin/pkill -SIGHUP -f midihub.py"
