# 44-midibt.rules
ACTION=="add|remove", SUBSYSTEM=="bluetooth", RUN+="/usr/local/bin/midihub.py"
