# 44-midibt.rules
ACTION=="add|remove", SUBSYSTEM=="sound", KERNEL=="card*", \
    RUN+="/usr/bin/pkill -SIGHUP -f midihub.py"
