# services/11-hookup.rules
ACTION=="add", SUBSYSTEM=="sound", RUN+="/bin/systemctl start hookup.service"
ACTION=="remove", SUBSYSTEM=="sound", RUN+="/bin/systemctl start hookup.service"
